# python/traning_cli/server/deploy/traning-garmin.service
[Unit]
Description=tRäning Garmin fetch (oneshot)
After=network-online.target

[Service]
Type=oneshot
User=krisse
EnvironmentFile=/etc/traning/env
WorkingDirectory=/home/krisse/dev/traning
ExecStart=/home/krisse/dev/traning/scripts/garmin_fetch_import.sh
StandardOutput=journal
StandardError=journal
