# python/traning_cli/server/deploy/traning-garmin.service
[Unit]
Description=tRäning Garmin fetch (oneshot)
After=network-online.target

[Service]
Type=oneshot
User=krisniem
EnvironmentFile=/etc/traning/env
WorkingDirectory=/home/krisniem/dev/traning
ExecStart=/home/krisniem/dev/traning/python/.venv/bin/traning fetch garmin --login-method browser
StandardOutput=journal
StandardError=journal
