# python/traning_cli/server/deploy/traning-receiver.service
[Unit]
Description=tRäning health data receiver (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=krisse
EnvironmentFile=/etc/traning/env
WorkingDirectory=/home/krisse/dev/traning
ExecStart=/home/krisse/dev/traning/python/.venv/bin/traning serve --host 0.0.0.0 --port 8421
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
