# python/traning_cli/server/deploy/traning-garmin.timer
[Unit]
Description=Garmin fetch fallback (every 2h, 06-22)

[Timer]
OnCalendar=*-*-* 06..22/2:00:00
Persistent=true

[Install]
WantedBy=timers.target
